# cron-sheduled-task/root/tmp/crontab
# This will run the task every minute

*/1 * * * * /bin/task.py
